# .config/sv/syncthing/run
#!/bin/sh

export HOME=/home/user/
exec 2>&1
exec /usr/bin/syncthing
